# scripts/com.shisui.debate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.shisui.debate</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/nayutasato/my-genspark/.venv/bin/python</string>
        <string>/Users/nayutasato/my-genspark/app.py</string>
        <string>debate-autonomous</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/nayutasato/my-genspark</string>

    <!-- 深夜3時に起動(夜間修行=com.shisui.study の深夜2時とはずらし、qwen2.5:32bの
         同時実行によるリソース競合を避ける)。Macが完全スリープ中は発火しない場合がある -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/nayutasato/my-genspark/output/study/launchd_debate_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/nayutasato/my-genspark/output/study/launchd_debate_stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
